# deploy/systemd/agma-provider-check.service.example
[Unit]
Description=AGMA billable live-provider validation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=agma
Group=agma
WorkingDirectory=/opt/agma/current
EnvironmentFile=/etc/agma/runtime.env
ExecStart=/usr/bin/node /opt/agma/current/runtime/dist/validation/live-provider-check.js --confirm-billable --config /var/lib/agma/runtime/config.yml
TimeoutStartSec=10min
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=

# This unit is intentionally not installable or scheduled. Start it manually
# only when the selected Provider account may incur the documented test calls.
